# templates/gitea-backup.timer.j2
[Unit]
Description=Gitea backup

[Timer]
OnCalendar={{ gitea_systemd_timer_on_calendar }}
RandomizedDelaySec={{ gitea_systemd_timer_randomized_delay_sec }}

[Install]
WantedBy=timers.target
